# root/etc/services.d/DnsServer/run
#!/usr/bin/with-contenv bash
cd /app || exit

exec \
	s6-setuidgid abc dotnet DnsServerApp.dll /config